# actax.service
[Unit]
Description=ActaX Server Panel
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/actax
# Zugangsdaten & Port kommen aus data/actax.env (von install.sh erzeugt: ACTAX_USER/ACTAX_PASS/ACTAX_PORT)
EnvironmentFile=-/opt/actax/data/actax.env
ExecStart=/opt/actax/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8080 --workers 1
# TLS/HTTPS (optional): Zertifikat unter System -> Sicherheit -> SSL erzeugen,
# obige Zeile auskommentieren und stattdessen verwenden:
# ExecStart=/opt/actax/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8080 --workers 1 --ssl-keyfile /opt/actax/data/certs/<CN>.key --ssl-certfile /opt/actax/data/certs/<CN>.crt
# danach: systemctl daemon-reload && systemctl restart actax  (Aufruf dann via https://)
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
